# deploy/stream-api.service
[Unit]
Description=Stream API - High-Fidelity Audio Streaming Server
After=network.target

[Service]
Type=simple
User=streamapi
Group=streamapi
WorkingDirectory=/opt/stream-api
ExecStart=/opt/stream-api/stream-api
EnvironmentFile=/opt/stream-api/.env
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

StandardOutput=journal
StandardError=journal

ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/stream-api /var/lib/stream-api
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
